# tutorials/songbird-node-installation/services/service-template.service
# songbird-service-template.service
# This template requires the run-songbird.sh script
[Unit]
Description=Songbird Systemd Service
StartLimitIntervalSec=0

[Service]
ExecStart=/home/ubuntu/run-songbird.sh
Restart=always
RestartSec=1
WorkingDirectory=/home/ubuntu
Environment=GOPATH=/home/ubuntu/go
Environment=PATH=/usr/local/go/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games:/snap/bin

User=ubuntu
Type=forking

StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=SongbirdNode

[Install]
WantedBy=multi-user.target


